# src/crontab
0 0 * * *  python3 /app/run.py >>/tmp/out.log 2>/tmp/err.log
